# deploy/supervisor.conf
# ============================================================
# KadrGo — Supervisor Konfiguratsiyasi
# /etc/supervisor/conf.d/kadrgo.conf
# ============================================================

# ─── Queue Worker (asosiy) ───
# Redis queue dan ishlarni bajaradi
[program:kadrgo-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/kadrgo/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --memory=128
numprocs=2
autostart=true
autorestart=true
startsecs=0
startretries=3
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/kadrgo/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600

# ─── Laravel Horizon (ixtiyoriy — Horizon o'rnatilgan bo'lsa) ───
# Horizon o'zi worker boshqaradi, yuqoridagi kadrgo-worker o'rniga ishlatish mumkin
# [program:kadrgo-horizon]
# process_name=%(program_name)s
# command=php /var/www/kadrgo/artisan horizon
# autostart=true
# autorestart=true
# user=www-data
# redirect_stderr=true
# stdout_logfile=/var/www/kadrgo/storage/logs/horizon.log
# stdout_logfile_maxbytes=10MB
# stdout_logfile_backups=5
# stopwaitsecs=3600
